# services/research_operations/deploy/native/systemd/research-operations-job-worker.service
[Unit]
Description=Research Operations admitted research job worker
Documentation=file:/opt/research-platform/current/services/research_operations/deploy/native/README.md
Requires=research-operations-preflight.service research-operations-migrate.service
After=research-operations-preflight.service research-operations-migrate.service
PartOf=research-operations.target
StartLimitIntervalSec=180
StartLimitBurst=5

[Service]
Type=simple
User=research-ops
Group=research-ops
WorkingDirectory=/opt/research-platform/current
EnvironmentFile=/etc/research-ops/runtime.env
Environment=RESEARCH_OPS_DATABASE_ROLE=runtime
Environment=PATH=/opt/research-platform/current/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
LoadCredential=operated-execution.key:/etc/research-ops/secrets/operated-execution.key
ExecStart=/usr/bin/env /opt/research-platform/current/services/research_operations/scripts/runtime-entrypoint.py /opt/research-platform/current/.venv/bin/research-ops research-job-worker --worker-id research-job:systemd-1 --poll-interval 1 --admission-lease-seconds 60
Restart=on-failure
RestartSec=5s
KillSignal=SIGTERM
KillMode=mixed
TimeoutStopSec=135s
TimeoutStopFailureMode=kill
OOMPolicy=stop
UMask=0027
RuntimeDirectory=research-operations
RuntimeDirectoryMode=0700
RuntimeDirectoryPreserve=yes
NoNewPrivileges=true
CapabilityBoundingSet=
AmbientCapabilities=
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
RestrictNamespaces=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
ReadOnlyPaths=/srv/research/data
ReadWritePaths=/srv/research/data/_internal_web/manifests
ReadWritePaths=/srv/research/artifacts /srv/research/reports /srv/research/cache /srv/research/registry
TasksMax=128
LimitNOFILE=8192
MemoryMax=4G
CPUQuota=400%
StandardOutput=journal
StandardError=journal
SyslogIdentifier=research-operations-job-worker
