# services/research_operations/deploy/native/nginx/nginx.service.d/research-operations.conf
[Service]
# Gunicorn sockets are 0660 research-ops:research-ops in 0750 runtime
# directories. Grant only the reverse proxy supplemental access; do not expose
# TCP listener ports or make either socket world-accessible.
SupplementaryGroups=research-ops
